# static/audit-daemon.service
[Unit]
Description=Storage daemon
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/audit-daemon --socket /run/audit-daemon.sock

[Install]
WantedBy=multi-user.target